# extra/autostart/birdtray.desktop
[Desktop Entry]
Name=Birdtray
Type=Application
Exec=birdtray
Icon=birdtray