# extra/autostart/birdtray.desktop
[Desktop Entry]
Name=Thuderbird
Type=Application
Exec=thunderbird
Icon=ceri-email