# extra/autostart/deluge.desktop
[Desktop Entry]
Name=Deluge
Type=Application
Exec=deluge
Icon=ceri-deluge