# extra/autostart/discord.desktop
[Desktop Entry]
Name=Discord
Type=Application
Exec=vesktop --enable-blink-features=MiddleClickAutoscroll --disable-smooth-scrolling
Icon=ceri-cord