# extra/autostart/dropbox.desktop
[Desktop Entry]
Name=Dropbox
Type=Application
Exec=dropbox
Icon=dropbox