# extra/autostart/floorp-work.desktop
[Desktop Entry]
Name=Floorp (Work)
Type=Application
Exec=floorp https://help.uillinois.edu/TDNext/Home/Desktop/Default.aspx