# extra/autostart/sound.desktop
[Desktop Entry]
Name=Startup sound
Type=Application
Exec=canberra-gtk-play -i desktop-login-first