# extra/autostart/steam.desktop
[Desktop Entry]
Name=Steam
Type=Application
Exec=steam -silent
Icon=ceri-steam