# extra/autostart/telegram.desktop
[Desktop Entry]
Name=Telgram
Type=Application
Exec=Telegram -startintray
Icon=org.telegram.desktop